# app/supervisord.conf
[supervisord]
nodaemon=true

[program:rnd]
command=bash /var/www/html/test.sh

[program:apache2]
command=service apache2 start

[program:stress]
command=stress -c 1